# systemd/one-piece-sniper-test-web.service
[Unit]
Description=One Piece Sniper web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Template only. Prefer scripts/install-user-services.sh to write concrete units.
Environment=ONE_PIECE_SNIPER_TEST_DIR=/ABSOLUTE/PATH/TO/one-piece-sniper-test
Environment=PORT=8787
WorkingDirectory=/ABSOLUTE/PATH/TO/one-piece-sniper-test
Environment=NODE_BIN=/usr/bin/node
Environment=PATH=/usr/bin:/bin
ExecStart=/usr/bin/node /ABSOLUTE/PATH/TO/one-piece-sniper-test/scripts/serve-static.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
