# systemd/one-piece-sniper-test-refresh.timer
[Unit]
Description=Daily 5:45 AM refresh timer for One Piece Sniper

[Timer]
OnCalendar=*-*-* 05:45:00
Persistent=true
Unit=one-piece-sniper-test-refresh.service

[Install]
WantedBy=timers.target
