# systemd/one-piece-sniper-test-refresh.service
[Unit]
Description=Refresh One Piece Sniper data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Template only. Prefer scripts/install-user-services.sh to write concrete units.
Environment=ONE_PIECE_SNIPER_TEST_DIR=/ABSOLUTE/PATH/TO/one-piece-sniper-test
Environment=NODE_BIN=/usr/bin/node
Environment=PATH=/usr/bin:/bin
WorkingDirectory=/ABSOLUTE/PATH/TO/one-piece-sniper-test
ExecStart=/usr/bin/bash /ABSOLUTE/PATH/TO/one-piece-sniper-test/scripts/refresh-data.sh
