# systemd/one-piece-sniper-test-refresh.timer
[Unit]
Description=Daily refresh timer for One Piece Sniper Test

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d
Persistent=true
Unit=one-piece-sniper-test-refresh.service

[Install]
WantedBy=timers.target
